# Linux/shellscript/deploy_info/long_script.conf
[program:long_script]
command=home/vinay/WorkSpace/Vin/Infra-Dev/Linux/shellscript/long.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
